# dotfiles-sync.timer
[Unit]
Description=Sync dotfiles every 2 weeks

[Timer]
OnBootSec=5min
OnUnitActiveSec=2w

[Install]
WantedBy=timers.target
