# dotfiles-sync.service
[Unit]
Description=Sync dotfiles to GitHub

[Service]
Type=oneshot
ExecStart=/home/vita_bite/dotfiles/sync-dotfiles.sh
